# scripts/quant-bot.service
[Unit]
Description=Quant Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/mnt/data/quant
Environment="PATH=/mnt/data/quant/.venv/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/mnt/data/quant"
EnvironmentFile=/mnt/data/quant/.env
ExecStart=/mnt/data/quant/.venv/bin/python3 -u src/scheduler/main.py
Restart=on-failure
RestartSec=30
StartLimitBurst=3
StartLimitIntervalSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quant-bot

[Install]
WantedBy=multi-user.target
